# .setup/submitty_grading_scheduler.conf
# launch the grading scheduler
start on filesystem
exec su -c '/usr/local/submitty/bin/grade_scheduler.py' hwcron