# .setup/submitty_user_scripts/submitty_daemon/crontab
* * * * * python3 ${SUBMITTY_INSTALL_DIR}/sbin/send_email.py
